# root/etc/s6-overlay/s6-rc.d/adguardhome-fixattrs/run
#!/command/with-contenv sh

chown -R root:appuser /opt/AdGuardHome/conf
chmod -R g+w /opt/AdGuardHome/conf
chown -R root:appuser /opt/AdGuardHome/work
chmod -R g+w /opt/AdGuardHome/work
