# root/etc/s6-overlay/s6-rc.d/adguardhome-fixattrs/up
/etc/s6-overlay/s6-rc.d/adguardhome-fixattrs/run
