# root/etc/s6-overlay/s6-rc.d/adguardhome/run
#!/command/with-contenv sh

exec /opt/AdGuardHome/AdGuardHome \
       --no-check-update \
       -c /opt/AdGuardHome/conf/AdGuardHome.yaml \
       -w /opt/AdGuardHome/work
